# etc/local.d/10-local.stop
## by F.C.
##
## This is a good place to stop all things
## started in /etc/local.d/10-local.start
